# etc/s6-overlay/s6-rc.d/uwsgi/run
#!/command/execlineb -P
/opt/rcamp_venv/bin/uwsgi /opt/uwsgi.ini
